# kali-fs/var/lib/dpkg/info/kali-defaults.list
/.
/etc
/etc/dconf
/etc/dconf/db
/etc/dconf/db/local.d
/etc/dconf/db/local.d/kali-tweaks
/etc/dconf/profile
/etc/dconf/profile/user
/etc/firefox-esr
/etc/firefox-esr/kali.js
/etc/firefox-esr/pref
/etc/firefox-esr/profile
/etc/profile.d
/etc/profile.d/kali.sh
/etc/systemd
/etc/systemd/network
/etc/udev
/etc/udev/rules.d
/etc/xdg
/etc/xdg/autostart
/etc/xdg/autostart/kali-noautomount.desktop
/etc/xdg/xfce4
/etc/xdg/xfce4/terminal
/etc/xdg/xfce4/terminal/terminalrc
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/postgresql@.service.d
/lib/systemd/system/postgresql@.service.d/kali_postgresql.conf
/usr
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/user
/usr/lib/systemd/user/pulseaudio.service.d
/usr/lib/systemd/user/pulseaudio.service.d/kali_pulseaudio.conf
/usr/lib/systemd/user/pulseaudio.socket.d
/usr/lib/systemd/user/pulseaudio.socket.d/kali_pulseaudio.socket.conf
/usr/share
/usr/share/doc
/usr/share/doc/kali-defaults
/usr/share/doc/kali-defaults/changelog.gz
/usr/share/doc/kali-defaults/copyright
/usr/share/firefox-esr
/usr/share/firefox-esr/distribution
/usr/share/firefox-esr/distribution/distribution.ini
/usr/share/glib-2.0
/usr/share/glib-2.0/schemas
/usr/share/glib-2.0/schemas/20_kali-defaults.gschema.override
/usr/share/kali-defaults
/usr/share/kali-defaults/.bashrc
/usr/share/kali-defaults/postgresql_reduce_shared_buffers
/usr/share/kali-defaults/setup-noautomount
/usr/share/kali-defaults/web
/usr/share/kali-defaults/web/css
/usr/share/kali-defaults/web/css/animate.css
/usr/share/kali-defaults/web/css/bootstrap.min.css
/usr/share/kali-defaults/web/css/course-tables.css
/usr/share/kali-defaults/web/css/custom-fonts.css
/usr/share/kali-defaults/web/css/fancybox
/usr/share/kali-defaults/web/css/fancybox/blank.gif
/usr/share/kali-defaults/web/css/fancybox/fancybox_loading.gif
/usr/share/kali-defaults/web/css/fancybox/fancybox_overlay.png
/usr/share/kali-defaults/web/css/fancybox/fancybox_sprite.png
/usr/share/kali-defaults/web/css/fancybox/jquery.fancybox.css
/usr/share/kali-defaults/web/css/flexslider.css
/usr/share/kali-defaults/web/css/font-awesome.css
/usr/share/kali-defaults/web/css/font-awesome.old.css
/usr/share/kali-defaults/web/css/overwrite.css
/usr/share/kali-defaults/web/css/roboto-font.css
/usr/share/kali-defaults/web/css/style.css
/usr/share/kali-defaults/web/fonts
/usr/share/kali-defaults/web/fonts/fa-brands-400.woff2
/usr/share/kali-defaults/web/fonts/fa-regular-400.woff2
/usr/share/kali-defaults/web/fonts/fa-solid-900.woff2
/usr/share/kali-defaults/web/fonts/source-sans
/usr/share/kali-defaults/web/fonts/source-sans/6xK1dSBYKcSV-LCoeQqfX1RYOo3qPZ7nsDI.woff2
/usr/share/kali-defaults/web/fonts/source-sans/6xK1dSBYKcSV-LCoeQqfX1RYOo3qPZ7psDJT9g.woff2
/usr/share/kali-defaults/web/fonts/source-sans/6xK3dSBYKcSV-LCoeQqfX1RYOo3qNq7lqDY.woff2
/usr/share/kali-defaults/web/fonts/source-sans/6xK3dSBYKcSV-LCoeQqfX1RYOo3qOK7l.woff2
/usr/share/kali-defaults/web/fonts/source-sans/6xKwdSBYKcSV-LCoeQqfX1RYOo3qPZY4lCdi18Smxg.woff2
/usr/share/kali-defaults/web/fonts/source-sans/6xKwdSBYKcSV-LCoeQqfX1RYOo3qPZY4lCds18Q.woff2
/usr/share/kali-defaults/web/fonts/source-sans/6xKwdSBYKcSV-LCoeQqfX1RYOo3qPZZclSdi18Smxg.woff2
/usr/share/kali-defaults/web/fonts/source-sans/6xKwdSBYKcSV-LCoeQqfX1RYOo3qPZZclSds18Q.woff2
/usr/share/kali-defaults/web/fonts/source-sans/6xKwdSBYKcSV-LCoeQqfX1RYOo3qPZZklydi18Smxg.woff2
/usr/share/kali-defaults/web/fonts/source-sans/6xKwdSBYKcSV-LCoeQqfX1RYOo3qPZZklyds18Q.woff2
/usr/share/kali-defaults/web/fonts/source-sans/6xKydSBYKcSV-LCoeQqfX1RYOo3i54rwlxdu.woff2
/usr/share/kali-defaults/web/fonts/source-sans/6xKydSBYKcSV-LCoeQqfX1RYOo3i54rwmRduz8A.woff2
/usr/share/kali-defaults/web/fonts/source-sans/6xKydSBYKcSV-LCoeQqfX1RYOo3ig4vwlxdu.woff2
/usr/share/kali-defaults/web/fonts/source-sans/6xKydSBYKcSV-LCoeQqfX1RYOo3ig4vwmRduz8A.woff2
/usr/share/kali-defaults/web/fonts/source-sans/6xKydSBYKcSV-LCoeQqfX1RYOo3iu4nwlxdu.woff2
/usr/share/kali-defaults/web/fonts/source-sans/6xKydSBYKcSV-LCoeQqfX1RYOo3iu4nwmRduz8A.woff2
/usr/share/kali-defaults/web/fonts/source-sans/source-sans.css
/usr/share/kali-defaults/web/homepage.html
/usr/share/kali-defaults/web/img
/usr/share/kali-defaults/web/img/slides
/usr/share/kali-defaults/web/img/slides/2.jpg
/usr/share/kali-defaults/web/js
/usr/share/kali-defaults/web/js/bootstrap.min.js
/usr/share/kali-defaults/web/js/jquery.easing.1.3.js
/usr/share/kali-defaults/web/js/jquery.js
/usr/share/kali-defaults/web/kali-site-logo.png
/usr/share/kali-defaults/web/offsec-man-in-the-door.png
/usr/share/kali-defaults/web/skins
/usr/share/kali-defaults/web/skins/default.css
/usr/share/kali-defaults/xsettings.xml
/etc/systemd/network/99-default.link
/etc/udev/rules.d/73-usb-net-by-mac.rules
